# svc/snapshot-monthly.timer
[Unit]
Description=create a monthly snapshot

[Timer]
OnCalendar=*-*-* *:00
Persistent=true

[Install]
WantedBy=timers.target
